# munin/files/Arch/munin-cron.service
[Unit]
Description=Munin Cron
After=network.target

[Service]
User=munin
ExecStart=/usr/bin/munin-cron
