# munin/files/Arch/munin-cron.timer
[Unit]
Description=Munin Cron Timer

[Timer]
OnCalendar=*-*-* *:00/5:00

[Install]
WantedBy=multi-user.target
